# components/ee/agent-smith/scripts/workspace.automount
[Unit]
Description=Workspace automount
ConditionPathExists=/workspace

[Automount]
Where=/workspace
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
